# build/files/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:geth]
command=/usr/local/bin/geth --datadir /root/.ethereum/ethchain-geth --http --http.addr 0.0.0.0 --http.corsdomain "*" --http.vhosts "*" --ws --ws.origins "*" --ws.addr 0.0.0.0 --rpc.gascap 0 --cache 4096 %(ENV_EXTRA_OPTS)s
autostart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0


[program:nginx]
command=nginx -c /etc/nginx/nginx.conf
autostart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
